# launchd/com.eugene.telegram_bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
"http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.eugene.telegram_bot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/eugene/Agents/supervisoragent/start_telegram_bot.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/eugene/Agents/supervisoragent</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/eugene/Agents/telegram_bot.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/eugene/Agents/telegram_bot.error.log</string>
</dict>
</plist>
